# systemd/profiled.service
[Unit]
Description=Profile Daemon

[Service]
Type=dbus
BusName=com.nokia.profiled
ExecStart=/usr/bin/profiled
